# scripts/systemd/english-news-digest-comments.timer
[Unit]
Description=Hourly English News Digest comment fetch check

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
Unit=english-news-digest-comments.service

[Install]
WantedBy=timers.target
